# AI100/Power_cycle.service
[Unit]
Description=Power Cycle test

[Service]
Type=oneshot
ExecStart=/usr/bin/reboot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
